# software/buildroot-2022.02.2-loshark/user-overlay/etc/profile.d/common.sh
export HOME=/root
export TERM=linux
export PS1='\u@\h:\w\$ '

ulimit -s 768
ulimit -l 128
