# software/buildroot-2022.02.2-loshark/user-overlay/etc/profile.d/shortcut.sh
alias l='ls -ahl'
alias ll='ls -ahl'
alias la='ls -ahl'

